# debian/soapysdr-server.prerm
#!/bin/sh
set -e

if [ "$1" = "remove" ]; then
	systemctl stop SoapySDRServer
	systemctl disable SoapySDRServer
fi

#DEBHELPER#
