# system/SoapySDRServer.service.in
[Unit]
Description=SoapyRemote network server
Wants=network-online.target
After=network-online.target

[Service]
ExecStart=@CMAKE_INSTALL_FULL_BINDIR@/SoapySDRServer --bind
KillMode=process
Restart=on-failure
LimitRTPRIO=99

[Install]
WantedBy=multi-user.target
